# deploy/betterclaw.service
[Unit]
Description=BetterClaw AI Assistant
After=cloud-sql-proxy.service docker.service
Requires=cloud-sql-proxy.service

[Service]
Type=simple
ExecStartPre=/usr/bin/docker pull us-central1-docker.pkg.dev/betterclaw-prod/betterclaw/agent:latest
ExecStart=/usr/bin/docker run --rm \
  --name betterclaw \
  --env-file /opt/betterclaw/.env \
  --network=host \
  us-central1-docker.pkg.dev/betterclaw-prod/betterclaw/agent:latest \
  --no-onboard
ExecStop=/usr/bin/docker stop betterclaw
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
